# packaging/ubuntu-jammy/postinst
#!/bin/bash
set -e
export DEBIAN_FRONTEND=noninteractive

echo "[INFO] Reloading systemd daemon..."
systemctl daemon-reexec
systemctl daemon-reload

echo "[INFO] Enabling services..."
# Fixed service name: added missing dash in "houston-zfs-storage-alert"
systemctl enable houston-dbus.service || true
systemctl enable --now houston-zfs-storage-alert.timer || true
systemctl enable houston-zfs-storage-alert.service || true

echo "[INFO] Starting services..."
# Restarting services safely, allowing them to fail without aborting script
systemctl restart zed || true
systemctl restart houston-dbus.service || true

# Optional: Python dependencies (disabled by default)
echo "[INFO] Installing Python dependencies..."
pip3 install --upgrade pip
pip3 install --upgrade google-auth google-auth-oauthlib

echo "[SUCCESS] Services enabled, started, and system is ready."
