# packaging/debian/DEBIAN/postinst
#!/bin/bash
set -e

LOGFILE="/var/log/cockpit-zfs-install.log"

echo "[INFO] Installing missing dependencies..." | tee -a "$LOGFILE"
export DEBIAN_FRONTEND=noninteractive
apt-get update >> "$LOGFILE" 2>&1
apt-get install -y python3-libzfs python3-pip sqlite3 python3-fastapi python3-uvicorn >> "$LOGFILE" 2>&1

echo "[INFO] Reloading systemd..." | tee -a "$LOGFILE"
systemctl daemon-reload

# Enable services if they exist
for service in houston-dbus.service fastapi-notifications.service; do
    if systemctl list-unit-files | grep -q "$service"; then
        echo "[INFO] Enabling $service..." | tee -a "$LOGFILE"
        systemctl enable "$service" >> "$LOGFILE" 2>&1

        echo "[INFO] Starting $service..." | tee -a "$LOGFILE"
        systemctl start "$service" || echo "[WARN] Failed to start $service" | tee -a "$LOGFILE"
    else
        echo "[WARN] $service not found, skipping." | tee -a "$LOGFILE"
    fi
done

echo "[INFO] Installation completed successfully!" | tee -a "$LOGFILE"
exit 0
